# systemd/pantalla-kiosk@.service
[Unit]
Description=Pantalla_reloj Kiosk (Epiphany) for user %i
After=pantalla-xorg.service pantalla-openbox@%i.service
Requires=pantalla-openbox@%i.service
PartOf=pantalla-openbox@%i.service

[Service]
Type=simple
User=%i
Environment=DISPLAY=:0
ExecStart=/usr/local/bin/kiosk-epiphany
Restart=on-failure
RestartSec=3
StartLimitIntervalSec=0

[Install]
WantedBy=graphical.target
