# systemd/pantalla-kiosk-autorefresh@.path
[Unit]
Description=Procesa solicitudes de refresco de kiosk (%i)

[Path]
PathChanged=/var/lib/pantalla-reloj/state/kiosk-refresh.flag
Unit=pantalla-kiosk-autorefresh@%i.service

[Install]
WantedBy=multi-user.target

