# systemd/pantalla-xorg.service
[Unit]
Description=Pantalla_reloj Xorg server (:0, vt7)
After=systemd-user-sessions.service

[Service]
Type=simple
User=root
ExecStartPre=/usr/lib/pantalla-reloj/xorg-launch.sh --prepare-only
ExecStartPre=/bin/sh -c 'if [ -f /var/lib/pantalla-reloj/.Xauthority ]; then \\
  owner=$(stat -c %u /var/lib/pantalla-reloj/.Xauthority) && \\
  group=$(stat -c %g /var/lib/pantalla-reloj/.Xauthority) && \\
  install -d -m 0755 -o "$owner" -g "$group" /var/log/pantalla; \\
else \\
  install -d -m 0755 /var/log/pantalla; \\
fi'
ExecStart=/usr/lib/xorg/Xorg :0 -verbose 3 -nolisten tcp -background none vt7 -auth /var/lib/pantalla-reloj/.Xauthority -logfile /var/log/pantalla/xorg.log
Restart=on-failure

[Install]
WantedBy=multi-user.target
