# systemd/pantalla-kiosk-chrome@.service
[Unit]
Description=Pantalla_reloj Chrome Kiosk for user %i
After=pantalla-xorg.service pantalla-openbox@%i.service pantalla-dash-backend@%i.service network-online.target
Requires=pantalla-openbox@%i.service pantalla-dash-backend@%i.service
Wants=pantalla-xorg.service network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
User=%i
Group=%i
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority
Environment=GDK_BACKEND=x11
Environment=GTK_USE_PORTAL=0
Environment=GIO_USE_PORTALS=0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
EnvironmentFile=-/var/lib/pantalla-reloj/state/kiosk.env
ExecStartPre=/usr/local/bin/pantalla-kiosk-prestart
LogsDirectory=pantalla
LogsDirectoryMode=0755
# NOTA: El perfil del navegador kiosk (/var/lib/pantalla-reloj/state/chromium-kiosk) se crea
# exclusivamente en scripts/install.sh. El wrapper pantalla-kiosk solo verifica que existe.
# Si falta, el wrapper fallará con un mensaje claro indicando que hay que ejecutar install.sh.
ExecStart=/usr/local/bin/pantalla-kiosk
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

