# systemd/pantalla-openbox@.service
[Unit]
Description=Pantalla_reloj Openbox session for %i (DISPLAY=:0)
After=pantalla-xorg.service
Wants=pantalla-xorg.service

[Service]
User=%i
Environment=DISPLAY=:0
Environment=XAUTHORITY=/var/lib/pantalla-reloj/.Xauthority
ExecStartPre=/usr/bin/install -d -m 0700 -o %i -g %i /run/user/1000
ExecStartPre=/bin/sh -c 'test -r "$XAUTHORITY"'
ExecStart=/usr/bin/openbox --startup /opt/pantalla/openbox/autostart
Restart=on-failure

[Install]
WantedBy=multi-user.target
