# systemd/pantalla-kiosk-chromium@.service
[Unit]
Description=Pantalla_reloj Kiosk (Chromium) for user %i
After=pantalla-xorg.service pantalla-openbox@%i.service pantalla-dash-backend@%i.service
Requires=pantalla-openbox@%i.service pantalla-dash-backend@%i.service
Wants=pantalla-xorg.service

[Service]
User=%i
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority
Environment=GDK_BACKEND=x11
Environment=GTK_USE_PORTAL=0
Environment=GIO_USE_PORTALS=0
Environment=CHROMIUM_SCALE=0.84
Environment=CHROMIUM_USER_DATA_DIR=/home/%i/.local/share/pantalla-reloj/chromium
Environment=CHROMIUM_CACHE_DIR=/home/%i/.cache/pantalla-reloj/chromium
EnvironmentFile=-/var/lib/pantalla-reloj/state/kiosk.env

ExecStartPre=/bin/sh -c 'if [ ! -r "$XAUTHORITY" ]; then if [ -f /var/lib/pantalla-reloj/.Xauthority ]; then cp -f /var/lib/pantalla-reloj/.Xauthority "$XAUTHORITY" && chown %i:%i "$XAUTHORITY" && chmod 600 "$XAUTHORITY"; fi; fi'
ExecStartPre=/bin/sh -c 'i=0; while [ $i -lt 10 ]; do if test -r "$XAUTHORITY"; then exit 0; fi; sleep 1; i=$((i+1)); done; echo "XAUTHORITY no es legible tras 10 segundos: $XAUTHORITY" >&2; exit 1'
ExecStartPre=/bin/sh -lc 'xset -dpms; xset s off; xset s noblank || true'
ExecStartPre=/bin/sh -lc 'pat="pantalla-kiosk|chrome.chromium|chromium-browser.Chromium-browser|chromium.Chromium"; wmctrl -lx 2>/dev/null | awk -v pat="$pat" '\''$3 ~ pat {print $1}'\'' | xargs -r -n1 wmctrl -ic || true'
ExecStartPre=/bin/sh -c 'curl -sf --max-time 5 http://127.0.0.1:8081/api/health >/dev/null || echo "WARN: Backend no responde, pero continuando (el script tolera backend ausente)" >&2 || true'

ExecStart=/usr/local/bin/pantalla-kiosk-chromium
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
