# systemd/pantalla-xorg.service
[Unit]
Description=Pantalla_reloj Xorg server (:0, vt7)
After=systemd-user-sessions.service systemd-logind.service
Conflicts=display-manager.service
Wants=graphical.target systemd-logind.service

[Service]
Environment=KIOSK_USER=dani
ExecStart=/usr/lib/pantalla-reloj/xorg-launch.sh
Restart=always
RestartSec=2s

[Install]
WantedBy=graphical.target
