# packaging/systemd/pantalla-dash-backend@.service
[Unit]
Description=Pantalla reloj FastAPI backend (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
# StateDirectory=pantalla-reloj crea /var/lib/pantalla-reloj con owner=User y 0755
StateDirectory=pantalla-reloj
# Si el grupo %i existe, usarlo; si no, usar el grupo primario del usuario
Group=%i
WorkingDirectory=/opt/pantalla-reloj/backend
Environment=PYTHONUNBUFFERED=1
Environment=PORT=8081
ExecStart=/usr/local/bin/pantalla-backend-launch
Restart=on-failure
RestartSec=2
TimeoutStartSec=30
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

