# system/pantalla-kiosk.service
[Unit]
Description=Pantalla Dash Chromium Kiosk
After=graphical.target network-online.target

[Service]
Type=simple
User=dashsvc
Group=dashsvc
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/dashsvc/.Xauthority
ExecStart=/usr/bin/chromium-browser \
  --kiosk http://localhost:8080 \
  --enable-features=VaapiVideoDecoder,CanvasOopRasterization,UseSkiaRenderer \
  --ignore-gpu-blocklist --enable-zero-copy --disable-gpu-driver-bug-workarounds \
  --disable-background-timer-throttling --disable-renderer-backgrounding \
  --noerrdialogs --disable-session-crashed-bubble --overscroll-history-navigation=0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
