# system/pantalla-bg-generate.service
[Unit]
Description=Pantalla Dash · Generación diaria de fondos futuristas
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/pantalla-dash/env
WorkingDirectory=/opt/dash/scripts
ExecStart=/usr/bin/python3 /opt/dash/scripts/generate_bg_daily.py
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/var/log/pantalla-dash/bg.log
StandardError=append:/var/log/pantalla-dash/bg.log
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
