# system/pantalla-ui.service
[Unit]
Description=Pantalla Chromium UI (kiosk)
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User={{UI_USER}}
Environment=DISPLAY=:0
Environment=PANTALLA_UI_URL=http://127.0.0.1:8080/
ExecStartPre=/usr/bin/bash -c '/usr/bin/pkill -f chromium || true'
ExecStart=/usr/local/bin/pantalla-ui-launch.sh
Restart=always
RestartSec=2

[Install]
WantedBy=graphical.target
