# system/pantalla-bg-generate.timer
[Unit]
Description=Timer fondos IA (dinámico por config.json)

[Timer]
OnBootSec=30s
OnUnitActiveSec=60min
Unit=pantalla-bg-generate.service

[Install]
WantedBy=timers.target
