# system/pantalla-bg-generate.timer
[Unit]
Description=Pantalla Dash · Programador de fondos futuristas diarios

[Timer]
OnBootSec=30s
OnCalendar=*-*-* 07:00:00
OnCalendar=*-*-* 12:00:00
OnCalendar=*-*-* 19:00:00
Persistent=true
AccuracySec=1min
RandomizedDelaySec=120
Unit=pantalla-bg-generate.service

[Install]
WantedBy=timers.target
