# system/pantalla-ui.service
[Unit]
Description=Pantalla Chromium UI (kiosk)
After=pantalla-openbox.service default.target
Wants=pantalla-openbox.service

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
Environment=PANTALLA_UI_URL=http://127.0.0.1/
ExecStartPre=/usr/bin/bash -lc 'pkill -f "^/snap/bin/chromium( |$)" || true'
ExecStart=/usr/local/bin/pantalla-ui-launch.sh
KillMode=process
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
