# system/pantalla-dash-backend.service
[Unit]
Description=Pantalla Dash Backend Service
After=network.target

[Service]
Type=simple
User=dashsvc
Group=dashsvc
WorkingDirectory=/opt/pantalla
Environment="PATH=/opt/pantalla/.venv/bin:/usr/bin"
ExecStart=/opt/pantalla/.venv/bin/uvicorn backend.app:app --host 127.0.0.1 --port 8787 --proxy-headers --workers 2
Restart=on-failure
RestartSec=5
Environment="PANTALLA_CONFIG_PATH=/etc/pantalla-dash/config.json"

[Install]
WantedBy=multi-user.target
