# system/pantalla-xorg@.service
[Unit]
Description=Pantalla - Xorg headless seat on %i
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=%i
# Detectar ruta Xorg en runtime: /usr/lib/xorg/Xorg o /usr/bin/Xorg
ExecStart=/bin/bash -lc 'XBIN=""; for p in /usr/lib/xorg/Xorg /usr/bin/Xorg; do [ -x "$p" ] && XBIN="$p" && break; done; [ -n "$XBIN" ] || { echo "Xorg bin not found"; exit 1; }; exec "$XBIN" :0 -nolisten tcp vt1'
Restart=always
RestartSec=2
Environment=HOME=/home/%i
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=multi-user.target
