# system/user/pantalla-openbox.service
[Unit]
Description=Pantalla Dash Openbox Session (with D-Bus)
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
ExecStart=/usr/bin/dbus-launch --exit-with-session /usr/bin/openbox-session
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
